# roles/supervisor/templates/laravel-worker.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php {{ LARAVEL_HOME }}/artisan queue:work {{ DB_CONNECTION }} --sleep=3 --tries=3 --daemon
autostart=true
autorestart=true
user={{ DB_USERNAME }}
numprocs=8
redirect_stderr=true
stdout_logfile={{ LARAVEL_HOME }}/storage/logs/worker.log
